# crontab
0 0,6,12,18 * * * bundle exec rake job:process_link_view_records
*/30 * * * * echo "cron is running!"